# labs/tn1/setup/oracle.crontab
# crontab for oracle user
# run a setup_baseline script check every 10 minutes
# */1 * * * * touch /home/oracle/workshops/cron.touch
#*/5 * * * * /home/oracle/labs/setup_baseline.sh >> /home/oracle/workshops/base.log 
#*/10 * * * * $HOME/setup/setup_baseline.sh > /dev/null 2>&1 
0,10,20,30,40,50 * * * * $HOME/setup/setup_baseline.sh > /dev/null 2>&1